# backend/dancee_workflow/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:restate-server]
command=/usr/local/bin/restate-server
environment=RESTATE_ROCKSDB_TOTAL_MEMORY_SIZE="256 MiB",RESTATE_ADMIN__QUERY_ENGINE__MEMORY_SIZE="64 MiB",RESTATE_BIFROST__RECORD_CACHE_MEMORY_SIZE="32 MiB"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app]
command=/app/scripts/register-deployment.sh
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
